# runsvdir/.config/service/pipewire/run
#!/bin/sh
export XDG_RUNTIME_DIR=/run/user/1000/
exec dbus-run-session -- pipewire
