# kmonad/dot-config/service/kmonad/run
#!/bin/sh
exec kmonad ~/.config/kmonad.kbd
