# runsvdir/.config/service/syncthing/run
#!/bin/sh
exec syncthing serve --no-browser
